# site/minecraft/files/minecraft.service
[Unit]
Description=Minecraft Server

Wants=network.target
